# airootfs/etc/systemd/system/dnsmasq.service.d/ordering.conf
[Unit]
After=wifi-ap-wait-ip.service
Wants=wifi-ap-wait-ip.service
